# ubuntu/provisioners/roles/ui/files/rc.local
#!/bin/sh -e

[ -f /etc/default/fetch-crl ] && . /etc/default/fetch-crl
/usr/sbin/fetch-crl -q -p 2 &

exit 0
-
